# files/after/after.1.service
[Service]
Type=oneshot
ExecStart=false
ExecStop=sleep 1
