# files/simple/simple-follow-up.service
[Unit]
Requires=simple-initial.service
After=simple-initial.service

[Service]
ExecStart=echo "My dependencies were activated successfully!"
