# files/requires/requires.1.service
[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=false
