# files/timer/timer.service
[Service]
ExecStart=echo "Hello world!"
