# files/fork/fork.2.service
[Unit]
Requires=fork.1.service
After=fork.1.service

[Service]
Type=oneshot
ExecStart=bash -c 'echo "Hello world" | nc -N localhost 18'
