# files/exec/exec-follow-up.service
[Unit]
Requires=exec-initial.service
After=exec-initial.service

[Service]
ExecStart=echo "My dependencies were activated successfully!"
