# files/exec/exec-initial.service
[Service]
# a service of exec type is considered active once execve() succeeds
Type=exec
ExecStart=/non/existent/binary
