# files/oneshot/oneshot.2.service
[Unit]
Requires=oneshot.1.service
After=oneshot.1.service

[Service]
ExecStart=echo "My dependencies were activated successfully!"
