# files/install/install.service
[Service]
ExecStart=echo "Hello world!"

[Install]
WantedBy=default.target
