# files/ctrl-alt-del.target
[Unit]
Description=A control-c halt target
Requires=halt.service
