# files/default.target
[Unit]
Description=A minimal default target
